# services/openrc/jinshu
#!/sbin/openrc-run

description="Jinshu Cognitive Graph-Vector Database Server Engine"

# Chemin vers le nouvel exécutable (à adapter selon ton installation, ex: /usr/local/bin/jinshu)
command="/usr/local/bin/jinshu"
command_args="serve"
command_background="yes"

directory="/tmp"
pidfile="/run/${RC_SVCNAME}.pid"
retry="respawn"

depend() {
    need net
}
